# rc.local
#!/bin/sh -e
/fuse/blobmount.bash
exit 0
